# droplet/systemd/orchon-ws.service
[Unit]
Description=ORCHON WebSocket Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/orchon/droplet/ws
Environment="HOME=/root"
Environment="WS_PORT=8405"
ExecStart=/bin/bash -c 'source /root/orchon/.env && exec /usr/bin/node /root/orchon/droplet/ws/server.js'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
